# cron/systemd/security-audit.timer
[Unit]
Description=Daily Security Audit Timer

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
