# dotfiles/dot_config/systemd/user/backup-projects.service
[Unit]
Description=Daily Projects Backup

[Service]
Type=simple
ExecStart=%h/.local/bin/backup_projects.sh
