# dotfiles/dot_config/systemd/user/dev-cleanup.timer
[Unit]
Description=Weekly AI/Dev Cache Cleanup Timer

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
