# dotfiles/dot_config/systemd/user/dev-cleanup.service
[Unit]
Description=AI/Dev Cache Cleanup

[Service]
Type=simple
ExecStart=%h/.local/bin/dev_cleanup.sh
