# dotfiles/dot_config/systemd/user/dwall.timer
[Unit]
Description=Run Dynamic Wallpaper Changer hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
